# services/logstash-forwarder.service.conf
[program:logstash-forwarder]
command=/opt/logstash-forwarder/logstash-forwarder -config /storage/config.json
stdout_events_enabled = true
stderr_events_enabled = true
